# upstart.conf
start on runlevel [234]
respawn
env RUST_LOG=warn
# todo: set user id
exec /usr/bin/verwalter

